# templates/blackbox-exporter.service.j2
[Unit]
Description=BlackBox Exporter

[Service]
ExecStart=/usr/local/bin/blackbox_exporter

[Install]
WantedBy=default.target
